# deploy/systemd/telegram-alerts.service
[Unit]
Description=Bot Telegram Alerts
After=network.target

[Service]
Type=simple
User=bot
EnvironmentFile=/etc/bot/.env
ExecStart=/opt/bot/bin/telegram-alerts
Restart=always

[Install]
WantedBy=multi-user.target
